# agent/com.code.remote-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.code.remote-agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/drew/code/code-remote/agent/run.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/drew/code/code-remote/agent/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/drew/code/code-remote/agent/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/drew/code/code-remote/agent</string>
</dict>
</plist>
